# node.service
[Unit]
Description=Nodejs unit file

[Service]
ExecStart=/usr/bin/node /home/exam/index.js
